# systemd/als-dimmer.service.in
[Unit]
Description=ALS-Dimmer: Ambient Light Sensor Based Display Brightness Control
Documentation=https://github.com/hackboxguy/als-dimmer.git
After=network.target multi-user.target
Wants=network.target

[Service]
Type=simple
# Pre-initialization script to enable I2C device access
#ExecStartPre=/home/pi/micropanel/bin/init-hh983-fpdlink.sh
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/als-dimmer --config @CMAKE_INSTALL_PREFIX@/etc/als-dimmer/config.json
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

# Run as root for I2C device access
User=root

# Security hardening (relaxed for I2C/DDC device access)
# Note: Some systemd security directives can cause NAMESPACE errors (exit code 226)
# on certain systems (especially older kernels or certain systemd versions).
# All security directives below are commented out for maximum compatibility.
# Uncomment them individually to test which ones your system supports.
#NoNewPrivileges=true
#ReadWritePaths=/var/lib/als-dimmer /tmp
#PrivateTmp=true
#ProtectHome=true
#ProtectSystem=strict
#DeviceAllow=/dev/i2c-1 rw
#DeviceAllow=/dev/i2c-0 rw
#DevicePolicy=closed

# Resource limits
LimitNOFILE=1024
MemoryMax=50M
CPUQuota=10%

[Install]
WantedBy=multi-user.target
